# deploy/systemd/kjorebok-web.service
[Unit]
Description=Kjorebok Web
After=network.target

[Service]
Type=simple
User=nguyen
WorkingDirectory=/home/nguyen/dev/kjorebok
Environment=NODE_ENV=production
Environment=PORT=3002
EnvironmentFile=-/etc/kjorebok/web.env
ExecStart=/usr/bin/env pnpm --filter @kjorebok/web start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
